# emanuel-server-supervisor.conf
[program:emanuel_server]
command=/home/rafi/start_emanuel_server.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/emanuel_server.out.log
stderr_logfile=/var/log/emanuel_server.err.log
